# 2021/medfusion/pump/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:telnetd]
command=xinetd -dontfork -stayalive
autorestart=true

