# 2018/al_berry/RunningOnRPi/rc.local
# add in /etc/rc.local:

sudo -b -u alberry "/home/alberry/Ph0wn2018/alberry.sh"
